# x11vnc.service
[Unit]
Description=VNC Server for X11 on Steam Deck
After=graphical-session.target

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStart=/usr/bin/x11vnc -forever -shared -rfbport 5900 -listen 0.0.0.0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
